# github-ssh-sync.service
[Unit]
Description=Sync GitHub SSH public keys to authorized_keys
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'source /usr/share/nvm/init-nvm.sh && node /opt/github-ssh-sync/sync-github-keys.js'

[Install]
WantedBy=multi-user.target
